# escos/src/mmitss/EquippedVehicleTrajectoryAware/mmitss_trajectory_aware.service
[Unit]
Description=MMITSS trajectory aware
Requisite=siedb.service ieee1609.service mmitssctl.service
ConditionPathExists=/run/services/mmitss_trajectory_aware.run
After=siedb.service ieee1609.service mmitssctl.service
StartLimitInterval=1min
StartLimitBurst=2

[Service]
Type=simple
Slice=system-comm-app.slice
ExecStart=/opt/bin/mmitss_trajectory_aware
WatchdogSec=0
Restart=no
RestartSec=30

[Install]
WantedBy=sie-its-wave.target
